# config/supervisor-lemp.conf
[program:nginx]
command=nginx

[program:php-fpm7]
command=php-fpm7 --nodaemonize